# systemd/.config/systemd/user/battery-notify.service
[Unit]
Description=Battery state and low-battery notifications
After=graphical-session.target

[Service]
ExecStart=%h/.config/hypr/scripts/battery.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
