# systemd/.config/systemd/user/media-notify.service
[Unit]
Description=MPRIS media change notifications
After=graphical-session.target

[Service]
ExecStart=%h/.config/hypr/scripts/media.sh
Restart=on-failure

[Install]
WantedBy=graphical-session.target
